# etc/s6/xrdp/run
#!/bin/sh

rendergrp=$(stat -c %G /dev/dri/renderD128)
usermod -a -G $rendergrp dosbox

exec  /usr/sbin/xrdp -nodaemon
